# scripts/configs/upstart/all-redises.conf
description "run pyres and resque redis instances"

start on runlevel [23]

task

script
  for conf in `ls /etc/redis/`;
  do
    start redis CONF=$conf
  done
end script
